# scripts/systemd/breachlab-audit-alarm.service
[Unit]
Description=BreachLab out-of-band DB-write detector
Wants=breachlab-audit-alarm.timer
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/breachlab-platform
# Pure-shell detector against the db compose container. Source the
# .env so AUDIT_ALARM_DISCORD_WEBHOOK is available to the script if set.
EnvironmentFile=-/opt/breachlab-platform/.env
ExecStart=/opt/breachlab-platform/scripts/check-out-of-band-writes.sh
User=root
NoNewPrivileges=true
ProtectSystem=strict
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
