# scripts/breachlab-invariants.service
[Unit]
Description=BreachLab invariant checker (drift detection across phantom + ghost + specter)
Documentation=man:systemd.unit(5)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/breachlab-invariants.sh check-all
TimeoutStartSec=180
Nice=15
IOSchedulingClass=idle
EnvironmentFile=-/etc/breachlab-invariants.env
